# config/supervisor/filebeat.conf
[program:filebeat]
command=/usr/bin/filebeat -c /etc/filebeat/filebeat.yml -e
autostart=true
autorestart=true
stderr_logfile=/var/log/filebeat.log
stdout_logfile=/var/log/filebeat.log
